# deploy/debian/shelter.upstart
# Upstart job configuration file
#
# Copyright 2014 Rafael Dantas Justo. All rights reserved.
# Use of this source code is governed by a GPL
# license that can be found in the LICENSE file.

description   "Shelter system"
author        "Rafael Dantas Justo"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 5 10

script
  exec /usr/shelter/bin/shelter \
    --config=/usr/shelter/etc/shelter.conf
end script